# ansible/roles/pre-any/templates/rkt-gc.service.j2
[Unit]
Description=Garbage Collection for rkt

[Service]
Environment=GRACE_PERIOD=30m0s
Type=oneshot
ExecStart=/usr/bin/rkt gc --grace-period=${GRACE_PERIOD}
